# DEBIAN/postinst
#!/bin/bash

# 安装依赖配置
VENV_PATH="/opt/oled-display/venv"
DEPS_DIR="/usr/share/oled-display/deps"
CONFIG_FILE="/etc/oled-display.conf"

# 启用I2C接口（兼容新旧Raspberry Pi OS版本）
CONFIG_FILE_PATH="/boot/firmware/config.txt"
[ -f "$CONFIG_FILE_PATH" ] || CONFIG_FILE_PATH="/boot/config.txt"

if ! grep -q '^dtparam=i2c_arm=on' "$CONFIG_FILE_PATH"; then
    echo "启用I2C接口..."
    echo "dtparam=i2c_arm=on" | tee -a "$CONFIG_FILE_PATH" >/dev/null
fi

# 加载I2C内核模块
if ! lsmod | grep -q i2c_dev; then
    modprobe i2c-dev
fi

# 创建虚拟环境（使用系统站点包）
echo "创建Python虚拟环境..."
python3 -m venv --system-site-packages "$VENV_PATH" || {
    echo "虚拟环境创建失败"
    exit 1
}

# 安装核心依赖
echo "安装Python依赖..."
"$VENV_PATH/bin/pip" install --no-index \
    --find-links="file://$DEPS_DIR" \
    --no-deps \
    luma.oled luma.core ntplib || {
    echo "依赖安装失败"
    exit 1
}

# 确保有写入权限
if [ ! -w /etc ]; then
    echo "错误：没有/etc目录写入权限" >&2
    exit 1
fi

# 创建配置文件
echo "正在创建配置文件..."
sudo tee "$CONFIG_FILE" > /dev/null <<EOL
[DISPLAY]
width = 128
height = 32
scroll_speed = 0.1

[TIME]
ntp_server = ntp.ntsc.ac.cn
timeout = 3
EOL

# 设置配置文件权限
sudo chmod 644 "$CONFIG_FILE"
sudo chown root:root "$CONFIG_FILE"

# 验证文件是否创建成功
if [ ! -f "$CONFIG_FILE" ]; then
    echo "错误：配置文件创建失败！" >&2
    exit 1
else
    echo "配置文件已创建：$CONFIG_FILE"
    cat "$CONFIG_FILE"
fi

# 启用服务
echo "设置系统服务..."
systemctl daemon-reload
systemctl enable oled-display.service

# 根据是否在终端运行决定是否交互
if [ -t 0 ]; then
    echo "---------------------------"
    echo "OLED显示服务安装完成"
    echo "当前配置："
    cat "$CONFIG_FILE"
    echo ""
    read -p "是否要修改屏幕尺寸？[y/N] " answer
    
    if [[ "$answer" =~ ^[Yy] ]]; then
        while true; do
            read -p "请输入屏幕尺寸（格式：128x32）：" screen_size
            if [[ "$screen_size" =~ ^[0-9]+x[0-9]+$ ]]; then
                width=$(echo "$screen_size" | cut -d'x' -f1)
                height=$(echo "$screen_size" | cut -d'x' -f2)
                
                # 使用sed保持原有配置其他内容
                sed -i "/^\[DISPLAY\]/,/^\[/ s/^width = .*/width = $width/" "$CONFIG_FILE"
                sed -i "/^\[DISPLAY\]/,/^\[/ s/^height = .*/height = $height/" "$CONFIG_FILE"
                
                systemctl restart oled-display.service
                echo "配置已更新"
                break
            else
                echo "格式错误，请重新输入"
            fi
        done
    fi
fi

# 无论是否交互都启动服务（如果未运行）
if ! systemctl is-active --quiet oled-display.service; then
    echo "启动OLED显示服务..."
    systemctl start oled-display.service
fi

exit 0
